# vbuntupkg/home/toor/.config/autostart/auto.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Change Wallpaper
Comment=Change Wallpaper
Exec=/opt/VbuntuServer/new_wallpaper
Terminal=false
